# launchd/com.codex.official-reg-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.codex.official-reg-monitor</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>-m</string>
    <string>official_reg_monitor.cli</string>
    <string>fetch</string>
    <string>--db</string>
    <string>__ROOT__/monitoring.sqlite</string>
    <string>--out</string>
    <string>__ROOT__/snapshots</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>4</integer>
    <key>Minute</key>
    <integer>20</integer>
  </dict>
  <key>WorkingDirectory</key>
  <string>__ROOT__</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONPATH</key>
    <string>__ROOT__/src</string>
  </dict>
  <key>StandardOutPath</key>
  <string>__ROOT__/logs/monitor.log</string>
  <key>StandardErrorPath</key>
  <string>__ROOT__/logs/monitor.err.log</string>
</dict>
</plist>
